# systemd/aoa-closeout-inbox.path
[Unit]
Description=Watch aoa-sdk closeout inbox for reviewed manifests

[Path]
PathExistsGlob=/srv/aoa-sdk/.aoa/closeout/inbox/*.json
PathChanged=/srv/aoa-sdk/.aoa/closeout/inbox
Unit=aoa-closeout-inbox.service

[Install]
WantedBy=default.target
