# mechanics/checkpoint/parts/reviewed-session-handoff-runner/closeout-inbox-user-units/aoa-closeout-inbox.service
[Unit]
Description=Process aoa-sdk reviewed session closeout inbox

[Service]
Type=oneshot
WorkingDirectory=/srv/AbyssOS/aoa-sdk
ExecStart=/usr/bin/env python3 /srv/AbyssOS/aoa-sdk/mechanics/checkpoint/parts/reviewed-session-handoff-runner/scripts/process_closeout_inbox.py --root /srv/AbyssOS/aoa-sdk --json
